# recipes-misc/basilisk-ii/files/systemd-units/basilisk-ii.service
[Unit]
Description=BasiliskII
After=multi-user-.target

[Service]
Type=simple
ExecStart=/bin/basilisk-ii.sh

[Install]
WantedBy=multi-user.target
